# etc/systemd/system/ctxcwalogd.service.d/logfix.conf
# Fix environment issue

[Service]
Environment="LD_PRELOAD=/opt/fakelog.so"
